# etc/services.d/99-crond/run
#!/bin/sh

. /etc/envvars

if [ -x $(which crond) ]; then
	exec 2>&1
	exec crond -f
fi
